# services.d/consul-template/run
#!/usr/bin/with-contenv sh

exec consul-template -config /consul-template/config.d
